# com.burk.crimson-monitor-wake.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.burk.crimson-monitor-wake</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/jonathan/.openclaw/workspace/crimson-mandate-agent/monitor_wake.py</string>
    </array>
    <key>StartInterval</key>
    <integer>300</integer>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/jonathan/.openclaw/workspace/crimson-mandate-agent/monitor_wake.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/jonathan/.openclaw/workspace/crimson-mandate-agent/monitor_wake.err.log</string>
</dict>
</plist>
